# data/ubuntu/.config/autostart/fcitx.desktop
[Desktop Entry]
Type=Application
Name=Fcitx
Exec=fcitx-autostart
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
